# airootfs/etc/profile.d/proxy.sh
export http_proxy=http://proxy.iiit.ac.in:8080
export https_proxy=http://proxy.iiit.ac.in:8080
export ftp_proxy=http://proxy.iiit.ac.in:8080
export rsync_proxy=http://proxy.iiit.ac.in:8080
export no_proxy="localhost,127.0.0.1,iiit.ac.in,.iiit.ac.in,iiit.net,.iiit.net,172.16.0.0/12,192.168.0.0/16,10.0.0.0/8"
